# IB_gateway/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 1280x800x24 -nolisten tcp
autorestart=true
priority=10

[program:fluxbox]
command=/usr/bin/fluxbox -display :1
environment=DISPLAY=":1"
autorestart=true
priority=20

[program:x11vnc]
command=/usr/bin/x11vnc -display :1 -rfbport 5900 -listen 127.0.0.1 -forever -shared -rfbauth /home/ibgateway/.vnc/passwd
autorestart=true
priority=30

[program:novnc]
command=/usr/bin/websockify --web=/usr/share/novnc 127.0.0.1:6080 127.0.0.1:5900
autorestart=true
priority=40

[program:ibgateway]
command=/app/start-ibgateway.sh
environment=DISPLAY=":1",HOME="/home/ibgateway"
autorestart=true
priority=50

[program:web]
command=gunicorn config.wsgi:application --bind 127.0.0.1:8001 --workers 2 --access-logfile - --forwarded-allow-ips=%(ENV_FORWARDED_ALLOW_IPS)s
environment=FORWARDED_ALLOW_IPS="*"
autorestart=true
priority=60

[program:broker]
command=python manage.py broker_worker
autorestart=true
priority=70

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart=true
priority=80

