# IB_gateway/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
childlogdir=/tmp

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 1280x800x24 -nolisten tcp
user=ibgateway
autorestart=unexpected
startretries=3
priority=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopsignal=TERM
stopasgroup=true
killasgroup=true

[program:fluxbox]
command=/app/wait-for-x.sh /usr/bin/fluxbox -display :1
user=ibgateway
environment=DISPLAY=":1",HOME="/home/ibgateway"
autorestart=unexpected
startretries=3
priority=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopsignal=TERM
stopasgroup=true
killasgroup=true

[program:x11vnc]
command=/app/wait-for-x.sh /usr/bin/x11vnc -display :1 -rfbport 5900 -listen 127.0.0.1 -forever -shared -rfbauth /home/ibgateway/.vnc/passwd
user=ibgateway
environment=DISPLAY=":1",HOME="/home/ibgateway"
autorestart=unexpected
startretries=3
priority=30
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopsignal=TERM
stopasgroup=true
killasgroup=true

[program:novnc]
command=/usr/bin/websockify --web=/usr/share/novnc 127.0.0.1:6080 127.0.0.1:5900
user=ibgateway
autorestart=unexpected
startretries=3
priority=40
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopsignal=TERM
stopasgroup=true
killasgroup=true

[program:ibgateway]
command=/app/wait-for-x.sh /app/start-ibgateway.sh
user=ibgateway
environment=DISPLAY=":1",HOME="/home/ibgateway"
autorestart=unexpected
startretries=5
startsecs=10
priority=50
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopsignal=TERM
stopasgroup=true
killasgroup=true

[program:web]
command=gunicorn config.wsgi:application --bind 127.0.0.1:8001 --workers 2 --access-logfile - --forwarded-allow-ips=%(ENV_FORWARDED_ALLOW_IPS)s
directory=/app
user=ibgateway
autorestart=unexpected
startretries=3
priority=60
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopsignal=TERM
stopasgroup=true
killasgroup=true

[program:broker]
command=python manage.py broker_worker
directory=/app
user=ibgateway
autorestart=unexpected
startretries=3
priority=70
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopsignal=TERM
stopasgroup=true
killasgroup=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart=unexpected
startretries=3
priority=80
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopsignal=QUIT
stopasgroup=true
killasgroup=true
