# Backend/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[program:web]
command=gunicorn config.wsgi:application --bind 0.0.0.0:%(ENV_PORT)s --workers 2 --access-logfile - --forwarded-allow-ips=%(ENV_FORWARDED_ALLOW_IPS)s
environment=PORT="8000",FORWARDED_ALLOW_IPS="*"
autorestart=true
priority=10

[program:celery]
command=celery -A config worker --loglevel=INFO --concurrency=2
autorestart=true
priority=20

[program:beat]
command=celery -A config beat --loglevel=INFO --schedule=/tmp/celerybeat-schedule
autorestart=true
priority=30

[program:market-consumer]
command=python manage.py consume_market_streams
autorestart=true
priority=25
