# Backend/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[program:web]
command=gunicorn config.asgi:application --worker-class uvicorn_worker.UvicornWorker --bind 0.0.0.0:%(ENV_PORT)s --workers 2 --access-logfile - --forwarded-allow-ips=%(ENV_FORWARDED_ALLOW_IPS)s
environment=PORT="8000",FORWARDED_ALLOW_IPS="*"
autorestart=true
priority=10

[program:celery-core]
command=celery -A config worker --loglevel=INFO --concurrency=2 --queues=celery,research_scoring,recommendation_cache,recommendations --hostname=core@%%h
autorestart=true
priority=20

[program:celery-data]
command=celery -A config worker --loglevel=INFO --pool=threads --concurrency=%(ENV_RESEARCH_MAX_PARALLEL_DATA_TASKS)s --queues=research_mapping,research_daily_data,research_intraday_data,research_fundamentals,research_events,research_features --hostname=data@%%h
autorestart=true
priority=21

[program:celery-research]
command=celery -A config worker --loglevel=INFO --pool=threads --concurrency=%(ENV_RESEARCH_MAX_PARALLEL_BACKTEST_TASKS)s --queues=research_single_asset,research_cross_sectional,research_allocators,research_overlays,research_pairs --hostname=research@%%h
autorestart=true
priority=22

[program:beat]
command=celery -A config beat --loglevel=INFO --schedule=/tmp/celerybeat-schedule
autorestart=true
priority=30

[program:market-consumer]
command=python manage.py consume_market_streams
autorestart=true
priority=25

[program:finnhub-market-data]
command=python manage.py consume_finnhub_market_data
autorestart=true
priority=26
